# deploy_tools/gunicorn.upstart.template.conf
description "Guinicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid ChaosOrd
chdir /home/ChaosOrd/sites/SITENAME/source

exec ../virtualenv/bin/gunicorn --bind unix:/tmp/SITENAME.socket superlists.wsgi:application

